# ssh-agent-replay/wsl-ssh-agent-relay.service
[Unit]
Description=WSL Relay SSH Agent

[Service]
Type=simple
ExecStart=/usr/local/sbin/wsl-ssh-agent-relay.sh

[Install]
WantedBy=default.target